# api-v3-lan-clients-ui.service.example
[Unit]
Description=APIv3 LAN Clients Console
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/api-v3-lan-clients-ui
Environment=HOST=0.0.0.0
Environment=PORT=8093
ExecStart=/opt/api-v3-lan-clients-ui/run_local.sh
Restart=always
RestartSec=3

# Change these if your Linux user is not ubuntu.
# The user/group below must be able to read and write /opt/api-v3-lan-clients-ui
# because run_local.sh creates/uses the local Python virtual environment.
User=ubuntu
Group=ubuntu

[Install]
WantedBy=multi-user.target
